# ansible/roles/backup_client/templates/borg-backup.timer.j2
[Unit]
Description=Borg Backup
[Timer]
OnUnitActiveSec={{ timer }}
[Install]
WantedBy=timers.target