# archive/roles/backup_client(Borg)/templates/borg-backup.timer.j2
[Unit]
Description=Borg Backup
[Timer]
{{ timer }}
[Install]
WantedBy=timers.target